# systemd/disk_alert.timer

[Unit]
Description=Run disk_alert daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
